# tools/stunnel.service.in
[Unit]
Description=SSL tunnel for network daemons
After=syslog.target

[Service]
ExecStart=@bindir@/stunnel
Type=forking

[Install]
WantedBy=multi-user.target
